# install-files/plugnmeet-recorder@main.service
[Unit]
Description=plugNmeet recorder main
Requires=plugnmeet-recorder.service
BindsTo=plugnmeet-recorder.service

[Service]
Type=idle
WorkingDirectory=/opt/plugNmeet/recorder
ExecStart=/usr/bin/node /opt/plugNmeet/recorder/dist/index.js
StandardOutput=syslog
StandardError=syslog
TimeoutStartSec=10
Restart=on-failure
RestartSec=10

[Install]
WantedBy=plugnmeet-recorder.service